# ops/systemd/channels-monitor-docker-prune.service
[Unit]
Description=Channels Monitor Docker build cache prune
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/channels-monitor-docker-prune.sh

